# systemd-services/pipewire-listen.service
[Unit]
Description=PipeWire Audio Event Listener - Auto-starts EQ Editor
After=pipewire.service
Wants=pipewire.service

[Service]
Type=simple
ExecStart=/bin/bash -c "pw-mon | grep --line-buffered 'changed' | while read -r line; do \
	STATUS=$(wpctl status); \
	BACKUP_LOCK=\"$HOME/.local/share/pipewire/backups/.backup_in_progress\"; \
	if echo \"$STATUS\" | grep 'Audio:playback' && echo \"$STATUS\" | grep 'Pipewire Gold'; then \
		systemctl --user start pipewire-eq.service; \
	else \
		if [[ ! -f \"$BACKUP_LOCK\" ]]; then \
			systemctl --user stop pipewire-eq.service; \
		fi; \
	fi; \
done"

Restart=always
RestartSec=5

# Resource limits
MemoryMax=32M
CPUQuota=5%
Nice=15

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pipewire-listen

[Install]
WantedBy=default.target
