# systemd-services/pipewire-listen.service
[Unit]
Description=PipeWire Audio Event Listener - Auto-starts EQ Editor
After=pipewire.service
Wants=pipewire.service

[Service]
Type=simple
ExecStart=/bin/bash -c "BACKUP_LOCK=\"$HOME/.local/share/goldaudio/backups/.backup_in_progress\"; pw-mon 2>/dev/null | grep --line-buffered -E 'Pipewire Gold.*Audio:playback' | while read -r line; do \
	STATUS=$(timeout 3 wpctl status 2>/dev/null || true); \
	if grep -qE 'Pipewire Gold.*Audio:playback' <<< "$STATUS"; then \
		systemctl --user is-active --quiet pipewire-eq.service || systemctl --user start pipewire-eq.service; \
	else \
		if [[ ! -f \"$BACKUP_LOCK\" ]]; then \
			systemctl --user is-active --quiet pipewire-eq.service && systemctl --user stop pipewire-eq.service; \
		fi; \
	fi; \
done"

Restart=always
RestartSec=5

# Resource limits
MemoryMax=32M
CPUQuota=3%
Nice=15

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=pipewire-listen

[Install]
WantedBy=default.target
